# autostart/auto.conf
# upstart service settings
description "cloudcare datakit"
start on runlevel [2345]
stop on runlevel [!2345]
env HOME=/root
export HOME
respawn
respawn limit 10 5
umask 022
exec /usr/local/cloudcare/forethought/datakit/datakit --config /usr/local/cloudcare/forethought/datakit/datakit.conf
post-stop exec sleep 1